# batterywarn.service
[Unit]
Description=Send a notification

[Service]
ExecStart=/usr/local/bin/battery.sh
